# scripts/dapps-updater.timer
[Unit]
Description=Periodic timer for DAPPS triggered-update worker
Documentation=https://github.com/packet-net/dapps

[Timer]
# Fire on a 60s cadence. Each tick the worker checks for the
# request marker dapps writes via POST /Update/apply; if absent,
# the worker exits in <100ms with "nothing to do". If present,
# it runs the full apply-update flow (download, swap, restart,
# 60s health window) and exits. Same shape as systemd-tmpfiles or
# any other periodic-cleanup timer.
OnBootSec=60s
OnUnitActiveSec=60s

# Persistent=true makes systemd remember the last successful run
# across reboots - if the host was off when the timer "should" have
# fired, we run it once on boot rather than skipping. Match the
# original C5.1 once-on-boot behaviour for the visibility check.
Persistent=true

[Install]
WantedBy=timers.target
