# scripts/dapps-updater.service
[Unit]
Description=DAPPS triggered-update worker (privileged)
Documentation=https://github.com/M0LTE/dapps
# We want to be sure dapps itself isn't actively transitioning when
# we start; ConditionPathExists guards against the timer firing on a
# host where dapps isn't installed at all.
ConditionPathExists=/opt/dapps/dapps

[Service]
# Type=oneshot — the .timer fires this on a 60s cadence; each
# invocation either no-ops (no marker file) or does one full update
# attempt (download, swap, restart, verify, possibly rollback).
Type=oneshot
User=root
Group=root
# The same dapps binary runs the update logic in a side-door mode
# that exits without booting the host. Plan C5.2 — keeps the update
# code in ONE place rather than duplicated in a shell script.
ExecStart=/opt/dapps/dapps --apply-update
WorkingDirectory=/var/lib/dapps

# This service is ONLY run by the timer; manual start is fine for a
# sysop wanting to force-poke the updater. Keep failure visibility on:
# a non-zero exit lands in the journal with the orchestrator's logs.

# Hardening — the updater needs write access to /opt/dapps and
# /var/lib/dapps and outbound HTTPS. Nothing else.
# NoNewPrivileges stays false — we ARE privileged on purpose.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/dapps /var/lib/dapps
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
